# amd-gpu-fan.service
[Unit]
Description="a simple daemon to control the fan curvature of amd video cards written in rust"
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
Environment=AMD_GPU_CONFIG_FILE=/etc/amd-gpu-fan/devices_config.json
ExecStart=/bin/amd-gpu-fan-daemon -rd

[Install]
WantedBy=multi-user.target
