# deb/DEBIAN/postinst
#!/bin/sh
systemctl enable amd-gpu-fan.service
systemctl start amd-gpu-fan.service
