# deb/DEBIAN/prerm
#!/bin/sh
systemctl stop amd-gpu-fan.service
systemctl disable amd-gpu-fan.service
